# telegram-hr-bot.service
[Unit]
Description=DEBUGGING Telegram HR Bot Environment

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/telegram-hr-bot

# Мы по-прежнему используем EnvironmentFile, чтобы проверить, читает ли он его
EnvironmentFile=/home/ubuntu/telegram-hr-bot/.env.production

# --- КЛЮЧЕВОЕ ИЗМЕНЕНИЕ: Запускаем ТОЛЬКО отладочный скрипт ---
ExecStart=/home/ubuntu/telegram-hr-bot/debug_env.sh

[Install]
WantedBy=multi-user.target
